# Scripts/rc.local
python3 /home/pi/gpiokb.py &
bash /home/pi/setspi.sh &
python3 /home/pi/GUI/Startup.py &
exit 0